# skel/.etc/profile.d/030-zprofile.zsh
# -*- mode: shell-script;-*-

##
#  zprofile bash compatibility
#

# /etc/zprofile -> /etc/profile
# emulate sh -c "source /etc/profile"

