# skel/.etc/profile.d/050-profile.sh
# -*- mode: shell-script;-*-

##
#  profile bash initialization
#


