# skel/.etc/profile.d/150-locale.sh
#!/bin/sh
# -*- mode: shell-script;-*-

##
#  locale definition
#

env_locale() {

LC_COLLATE=C
LANG=en_US.UTF-8

export LC_COLLATE
export LANG
    
}

#env_locale
