# skel/.etc/profile.d/150-locale.sh
# -*- mode: shell-script;-*-

##
#  locale definition
#

function env_locale {

export LC_COLLATE=C
export LANG=en_US.UTF-8
    
}

#env_locale
