# skel/.etc/profile.d/310-java.bash
# -*- mode: shell-script;-*-

##
#  Java environment 
#

#set -x

#set +x
