# skel/.etc/profile.d/330-ruby.sh
# -*- mode: shell-script;-*-

##
#  RVM ruby environment 
#

export SYS_RUBY=`which ruby`



