# skel/.etc/profile.d/330-ruby.sh
# -*- mode: shell-script;-*-

##
#  RVM ruby environment 
#

export SYS_RUBY=`which ruby`
export RVM_SCRIPT=~/.rvm/scripts/rvm

[ -x $RVM_SCRIPT ] && source $RVM_SCRIPT || true



