# skel/.etc/profile.d/330-ruby.sh
#!/bin/sh
# -*- mode: shell-script;-*-

##
#  RVM ruby environment 
#

SYS_RUBY=`which ruby`

export SYS_RUBY


