# skel/.etc/profile.d/370-haskell.zsh
# -*- mode: shell-script;-*-

##
#  Haskell cabal environment
#

