# skel/.etc/profile.d/510-include.sh
#!/bin/sh
# -*- mode: shell-script;-*-

##
#  application support
#

[ -x ~/.profile-share.sh ] && . ~/.profile-share.sh || true
[ -x ~/.profile-local.sh ] && . ~/.profile-local.sh || true


