# skel/.etc/profile.d/510-include.sh
# -*- mode: shell-script;-*-

##
#  application support
#

[ -x ~/.profile-share.sh ] && source ~/.profile-share.sh || true
[ -x ~/.profile-local.sh ] && source ~/.profile-local.sh || true


