# system/etc/profile.d/40cvs.sh
#!/bin/bash


if [ -d /apl/vcs/cvs/CVSROOT ]; then
	export  CVSROOT=:pserver:hute37@localhost:/apl/vcs/cvs
fi



