# system/etc/xinetd.d/cvspserver
service cvspserver
{
	disable = no
	socket_type = stream
	wait = no
	user = root
	group = root
	log_type = FILE /var/log/cvspserver
	protocol = tcp
	env  = '$HOME=/apl/vcs/cvs'
#	bind = 127.0.0.1
	log_on_failure += USERID
	port = 2401
	server = /usr/bin/cvs
	server_args = -f --allow-root=/apl/vcs/cvs pserver
}



